# s6/rclone/run
#!/bin/bash

set -e
set -x

/bin/mkdir -p /mnt/rclone

/usr/bin/rclone mount --allow-other -o uid=$MOUNT_UID -o gid=$MOUNT_GID $RCLONE_REMOTE: /mnt/rclone
