# s6/rclone/run
#!/bin/bash

set -e
set -x

/bin/mkdir -p /mnt/rclone

/usr/bin/rclone mount --allow-other $RCLONE_REMOTE: /mnt/rclone
