# s6/plexdrive/run
#!/bin/bash

set -e
set -x

/bin/mkdir -p /mnt/rcache

/usr/bin/rclone mount --allow-other --allow-non-empty --uid=$MOUNT_UID --gid=$MOUNT_GID $RCLONE_CACHE /mnt/rcache
