# s6/plexdrive/run
#!/bin/bash

set -e
set -x

/bin/mkdir -p /mnt/plexdrive

/usr/bin/plexdrive mount /mnt/plexdrive --uid=$MOUNT_UID --gid=$MOUNT_GID -o allow_other -o allow_non_empty_mount
